# bin/cron/com.aikb.monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ai-kb Phase 3 — monthly cron job. Runs 1st of month, 04:00 local time.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aikb.monthly</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VAULT__/bin/cron/monthly.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__VAULT__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__VAULT__/.kb/cron-monthly.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__VAULT__/.kb/cron-monthly.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
